# linux/.neocon/syntax_cache/inittab
[{'word': 'cSs', 'menu': '[S] inittabRunLevel'}, {'word': 'ctrlaltdel', 'menu': '[S] inittabActionNa'}, {'word': 'off', 'menu': '[S] inittabActionNa'}, {'word': 'end', 'menu': '[S] inittabShString'}, {'word': 'skip', 'menu': '[S] inittabShString'}, {'word': 'nextgroup', 'menu': '[S] inittabColonRun'}, {'word': 'xtgroup', 'menu': '[S] inittabColonRun'}, {'word': 'xxx', 'menu': '[S] inittabError'}, {'word': 'XXX', 'menu': '[S] inittabFixme'}, {'word': 'luster', 'menu': '[S] inittabSh'}, {'word': 'nsparent', 'menu': '[S] inittabProcess'}, {'word': 'start', 'menu': '[S] inittabShString'}, {'word': 'wait', 'menu': '[S] inittabActionNa'}, {'word': 'NOT', 'menu': '[S] inittabFixme'}, {'word': 'kbrequest', 'menu': '[S] inittabActionNa'}, {'word': 'sysinit', 'menu': '[S] inittabActionNa'}, {'word': 'bootwait', 'menu': '[S] inittabActionNa'}, {'word': 'boot', 'menu': '[S] inittabActionNa'}, {'word': 'ntains', 'menu': '[S] inittabAction'}, {'word': 'inittabError', 'menu': '[S] inittabError'}, {'word': 'inittabFixme', 'menu': '[S] inittabFixme'}, {'word': 'inittabComment', 'menu': '[S] inittabComment'}, {'word': 'inittabShString', 'menu': '[S] inittabShString'}, {'word': 'inittabShOption', 'menu': '[S] inittabShOption'}, {'word': 'inittabShCommand', 'menu': '[S] inittabShComman'}, {'word': 'inittabActionName', 'menu': '[S] inittabActionNa'}, {'word': 'initdefault', 'menu': '[S] inittabActionNa'}, {'word': 'inittabColonRunLevels', 'menu': '[S] inittabColonRun'}, {'word': 'inittabRunLevels', 'menu': '[S] inittabColonRun'}, {'word': 'inittabColonAction', 'menu': '[S] inittabColonRun'}, {'word': 'inittabId', 'menu': '[S] inittabId'}, {'word': 'inittabAction', 'menu': '[S] inittabColonAct'}, {'word': 'inittabColonProcess', 'menu': '[S] inittabAction'}, {'word': 'inittabProcessPlus', 'menu': '[S] inittabColonPro'}, {'word': 'inittabProcess', 'menu': '[S] inittabColonPro'}, {'word': 'inittabSh', 'menu': '[S] inittabProcess'}, {'word': 'powerokwait', 'menu': '[S] inittabActionNa'}, {'word': 'powerwait', 'menu': '[S] inittabActionNa'}, {'word': 'powerfailnow', 'menu': '[S] inittabActionNa'}, {'word': 'powerfail', 'menu': '[S] inittabActionNa'}, {'word': 'match', 'menu': '[S] inittabError'}, {'word': 'cluster', 'menu': '[S] inittabSh'}, {'word': 'tains', 'menu': '[S] inittabComment'}, {'word': 'FIXME', 'menu': '[S] inittabFixme'}, {'word': 'roup', 'menu': '[S] inittabId'}, {'word': 'alnum', 'menu': '[S] inittabShOption'}, {'word': 'TODO', 'menu': '[S] inittabFixme'}, {'word': 'transparent', 'menu': '[S] inittabProcess'}, {'word': 'once', 'menu': '[S] inittabActionNa'}, {'word': 'ondemand', 'menu': '[S] inittabActionNa'}, {'word': 'respawn', 'menu': '[S] inittabActionNa'}, {'word': 'contains', 'menu': '[S] inittabComment'}]
